# use-cases/elastic_router/opentsdb/supervisor-dd.conf
[program:ddproxy]
command=/opt/sei-bin/start_ddproxy.sh
priority=10
