# use-cases/elastic_router/opentsdb/supervisor-hbase.conf
[program:hbase]
command=/opt/sei-bin/start_hbase.sh
priority=20
stopwaitsecs=60
