# use-cases/elastic_router/opentsdb/supervisor-tsdb.conf
[program:opentsdb]
command=/opt/sei-bin/start_opentsdb.sh
priority=30
stopwaitsecs=60
