# NFs/docker/elastic-router/opentsdb/supervisor-dd.conf
[program:ddproxy]
command=/opt/sei-bin/start_ddproxy.sh